# assignement04/10-local.rules
ACTION=="add", SUBSYSTEM=="input", SUBSYSTEMS=="usb", RUN+="/sbin/insmod /transit/keyboard/hello.ko"
